# temp-sensor.service
[Unit]
Description=Temp Sensor Service
After=network.target
StartLimitIntervalSec=0
[Service]
WorkingDirectory=/home/imran/projects/temp_sensor
Type=simple
Restart=always
RestartSec=10
User=imran
ExecStart=/bin/bash /home/imran/projects/temp_sensor/startup.sh

[Install]
WantedBy=multi-user.target
